# orangepi/wlan/10-usbmount.rules
KERNEL=="sd*[!0-9]|sr*", ENV{ID_SERIAL}!="?*", SUBSYSTEMS=="usb", DRIVERS=="usb-storage", RUN+="/usr/bin/usbnetworkconfig.sh"

